# testdata/missing_service/orphan.timer
[Unit]
Description=Timer with no paired service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
