# testdata/orphan/stranded.service
[Unit]
Description=No matching timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stranded.sh

[Install]
WantedBy=timers.target
