# testdata/bad_calendar/broken.service
[Unit]
Description=Companion service for broken timer

[Service]
Type=oneshot
ExecStart=/bin/true
