# testdata/bad_calendar/broken.timer
[Unit]
Description=Timer with malformed OnCalendar

[Timer]
OnCalendar=*-*-* 99:99:99
Persistent=true

[Install]
WantedBy=timers.target
