# testdata/noisy/checkup.service
[Unit]
Description=Companion service for the every-minute timer

[Service]
ExecStart=/usr/local/bin/checkup.sh
