# testdata/conflicting/backup.service
[Unit]
Description=Backup A

[Service]
Type=oneshot
ExecStart=/var/backups/run.sh
