# testdata/conflicting/rotate.service
[Unit]
Description=Rotate B

[Service]
Type=oneshot
ExecStart=/var/backups/run.py
